# systemd/update-clean.timer
[Unit]
Description=Weekly system update and cleanup

[Timer]
OnCalendar=Sun *-*-* 04:00:00
Persistent=true
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target